# ai_coding.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aicoding.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhourongbing/Downloads/.venv/bin/python3</string>
        <string>/Users/zhourongbing/Downloads/Ai编程/AIcoding.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/zhourongbing/Downloads/Ai编程/logs/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/zhourongbing/Downloads/Ai编程/logs/output.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/zhourongbing/Downloads/Ai编程</string>
</dict>
</plist> 